# crontab
* * * * * root /usr/local/sousvide/reset-pin
